# templates/default/openvswitch_sdn.conf.erb
start on (filesystem and net-device-up)
stop on runlevel [016]

# Automatically restart process if crashed
respawn

pre-start script
end script

script
  exec <%= node[:openvswitch][:sdn_install_dir] %>/pox/pox.py inception \
    --ip_prefix=<%= node[:openvswitch][:ip_prefix] %> log.level \
    --DEBUG log --file=/tmp/pox.log
end script

pre-stop script
  exec ps -ef | grep -i pox.py | grep -i python | \
    awk '/root/{print $2}' | xargs kill -9
end script

post-stop script
end script
